# com.token-dashboard.aggregate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.token-dashboard.aggregate</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>${HOME}/Documents/dec/token-dashboard/scripts/aggregate.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>${HOME}/Documents/dec/token-dashboard</string>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>${HOME}/Documents/dec/token-dashboard/logs/aggregate.log</string>
    <key>StandardErrorPath</key>
    <string>${HOME}/Documents/dec/token-dashboard/logs/aggregate.error.log</string>
</dict>
</plist>
