# com.token-dashboard.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.token-dashboard.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>${HOME}/Documents/dec/token-dashboard/server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>${HOME}/Documents/dec/token-dashboard</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>${HOME}/Documents/dec/token-dashboard/logs/server.log</string>
    <key>StandardErrorPath</key>
    <string>${HOME}/Documents/dec/token-dashboard/logs/server.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>8901</string>
    </dict>
</dict>
</plist>
